# rpi-4b/rootfs_full/etc/inittab
#This is run first except when booting in single-user mode.
#use software uart,avoid //bin.sh:can't access tty: job control turned off(such as ctrl+c)
::sysinit:/etc/init.d/rcS
#/bin/sh invocations on delected ttys
::respawn:-/bin/sh
#Start an "askfirst" shell on the console(whatever that may be)
#use software uart,avoid //bin.sh:can't access tty: job control turned off(such as ctrl+c)
::askfirst:-/bin/sh
#Stuff to do when restartsing the init process
::restart:/sbin/init
#Stuff to do before rebooting
::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/swapoff -a
